# modules/cffex-daily/com.yuque.cffex-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.yuque.cffex-daily</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/mac/Desktop/github/my_tool_project/modules/cffex-daily/run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>21</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mac/Desktop/github/my_tool_project</string>
    <key>StandardOutPath</key>
    <string>/Users/mac/Desktop/github/my_tool_project/modules/cffex-daily/work/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mac/Desktop/github/my_tool_project/modules/cffex-daily/work/logs/launchd.err.log</string>
    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
